# etc/systemd/filkoll-update.service.in
[Unit]
Description=filkoll cache update
RequiresMountsFor="/var/cache/filkoll"
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart={BINDIR}/filkoll update

# Sandboxing
PrivateUsers=true
CapabilityBoundingSet=
LockPersonality=true
MemoryDenyWriteExecute=true
NoNewPrivileges=true
PrivateIPC=true
PrivateMounts=true
PrivateNetwork=true
PrivateTmp=true
ProcSubset=pid
ProtectClock=true
ProtectControlGroups=true
ProtectHome=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectProc=noaccess
ProtectSystem=strict
RestrictAddressFamilies=none
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native

# Allow what we need
ReadWritePaths=/var/cache/filkoll

# Logging
StandardOutput=journal
StandardError=journal

# Priority
Nice=19
CPUSchedulingPolicy=idle
IOSchedulingClass=idle
